# wifimngr/files/wifimngr.init
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

PROG=/usr/bin/ucode
WIFIMNGR_DIR=/usr/share/wifimngr

start_service() {
	procd_open_instance
	procd_set_param command "$PROG" -S -L "$WIFIMNGR_DIR" "$WIFIMNGR_DIR/wifimngr.uc"
	procd_set_param respawn
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "wireless"
}

reload_service() {
	stop
	start
}
